# restapi/restapi.conf
description "restapi"
start on runlevel [2345]
stop on runlevel [06]
respawn

chdir /home/gary/Documents/fyp_2021_busq/restapi
exec npm start