# upstart/data_pool.conf
description "data_pool"
start on runlevel [2345]
stop on runlevel [06]
respawn

chdir /home/gary/Documents/fyp_2021_busq
exec python3.5 data_pool.py