# upstart/gen_filter_list.conf
description "gen_filter_list"
start on runlevel [2345]
stop on runlevel [06]
respawn

chdir /home/gary/Documents/fyp_2021_busq
exec python3.5 gen_filter_list.py